# include/profile.d/grep.sh
export GREP_OPTIONS="--binary-files=without-match --color=auto"
