# include/profile.d/less.sh
export LESS="--chop-long-lines --hilite-search --ignore-case --LINE-NUMBERS --RAW-CONTROL-CHARS --quit-if-one-screen --status-column --tabs=4"
