# include/profile.d/mailcheck.sh
# To avoid delays you should either unset MAILCHECK or define MAILPATH to point to your existing mail inbox
unset MAILCHECK
